# backend/cubie-backend.service
[Unit]
Description=CubieCloud Backend API
After=network.target

[Service]
Type=simple
User=cubie
Group=cubie
WorkingDirectory=/opt/cubie/backend
ExecStart=/opt/cubie/backend/venv/bin/python -m app.main
Restart=always
RestartSec=5

# Environment overrides (edit as needed)
# NOTE: Do NOT hardcode `CUBIE_JWT_SECRET` here. Use one of the following safer options:
# - Provide `CUBIE_JWT_SECRET` as an environment variable from a secure system-wide secret manager.
# - Place the secret in `/var/lib/cubie/jwt_secret` (file is read by the app when no env var provided).
# - If you need to set environment variables from a file, use `EnvironmentFile=` with a file owned by root and mode 600.
# The service intentionally avoids embedding the JWT secret in this unit file to reduce accidental leaks.
Environment=CUBIE_DEVICE_SERIAL=CUBIE-A7A-2025-001
Environment=CUBIE_PAIRING_KEY=your-pairing-key
Environment=CUBIE_NAS_ROOT=/srv/nas
Environment=CUBIE_DATA_DIR=/var/lib/cubie
Environment=CUBIE_TLS_ENABLED=true
# Optional extra CORS origins (comma-separated). Leave empty to keep defaults:
# http://localhost,http://localhost:3000
Environment=CUBIE_CORS_ORIGINS=

[Install]
WantedBy=multi-user.target
